# docs/supervisor.conf
[program:jms]
command=/opt/py3/bin/python /opt/jumpserver/jms start all
user=root
numprocs=1
autostart=true
autorestart=true
startsecs=3
stdoutlogfile=/opt/jumpserver/logs/supervisor.log
